# conf/supervisor.conf.in
[program:kevro]
command={{ deploy_root }}/python/bin/gunicorn wsgi -w 4 --max-requests=10000 --bind 127.0.0.1:{{ port_prefix }}11
directory={{ deploy_root }}/{{ app_name }}/
user=www-data
stdout_logfile={{ deploy_root }}/log/{{ app_name }}-django.log
stderr_logfile={{ deploy_root }}/log/{{ app_name }}-django-error.log
stopsignal=QUIT
environment=
    PYTHON_EGG_CACHE='/var/www-data/.python-eggs',
    DJANGO_SETTINGS_MODULE=project.settings
